# debian/postrm
#! /bin/sh
# postrm script for magicolor5430dl
#
# see: dh_installdeb(1)

#set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove)
       PPD_DIR="/etc/cups/ppd/"
       FILETYPE=".ppd"
       LPADMIN="/usr/sbin/lpadmin"
       if [ -x $LPADMIN ] && [ -d $PPD_DIR ] && [ -n "`ls -1A ${PPD_DIR}`" ]; then
          KMPPDS=(`grep -l "magicolor 5430 DL" ${PPD_DIR}*`)
          for kmppd in ${KMPPDS[@]}; do
             ppdFile=${kmppd#${PPD_DIR}}
             Printer=${ppdFile%${FILETYPE}}
             $LPADMIN -x $Printer
             if [ $? -eq 0 ]; then
                echo Succeed to delete $Printer
             else
                echo Fail to delete $Printer
             fi
          done
        fi

	[ -x /etc/init.d/cupsys ] && invoke-rc.d cupsys force-reload || true
	;;

    purge|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
